# kiosk/systemd/muse-prepare.service
[Unit]
Description=Prepare Muse data and migrations before listeners start
After=local-fs.target
Before=muse-main.service muse-phone-upload.service

[Service]
Type=oneshot
User=muse
Group=muse
EnvironmentFile=/etc/muse/muse.env
EnvironmentFile=-/run/muse/network.env
ExecStart=/opt/muse/current/kiosk/muse-backend prepare-production
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX
SystemCallArchitectures=native
ReadWritePaths=/var/lib/muse /run/muse
RemainAfterExit=true
TimeoutStartSec=10min

[Install]
WantedBy=multi-user.target
